# packages/node-deploy/docker/gateway/supervisord.conf
# Chicago Forest Network Gateway - Supervisord Configuration
# Manages multiple services in the gateway container

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Forest Gateway Node
[program:forest-gateway]
command=/usr/bin/node /app/packages/node-deploy/dist/index.js --gateway
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/forest/gateway.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/forest/gateway-error.log
stderr_logfile_maxbytes=50MB
environment=NODE_ENV="production",GATEWAY_MODE="true"

# DNSMasq - DHCP and DNS
[program:dnsmasq]
command=/usr/sbin/dnsmasq -k --conf-file=/etc/dnsmasq.d/forest.conf
autostart=%(ENV_ENABLE_DNS)s
autorestart=true
startsecs=5
stdout_logfile=/var/log/forest/dnsmasq.log
stderr_logfile=/var/log/forest/dnsmasq-error.log

# WireGuard Interface Manager
[program:wireguard]
command=/bin/sh -c "wg-quick up wg-forest && sleep infinity"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/forest/wireguard.log
stderr_logfile=/var/log/forest/wireguard-error.log

# Health Check Server
[program:healthcheck]
command=/usr/bin/node /app/packages/node-deploy/dist/healthcheck.js
directory=/app
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/forest/healthcheck.log
stderr_logfile=/var/log/forest/healthcheck-error.log

[group:forest]
programs=forest-gateway,dnsmasq,wireguard,healthcheck
priority=999
